# crontab
* * * * * echo "Cron is working at $(date)" >> /var/log/cron.log 2>&1
